# python/iems5722_task.conf
[program:iems5722_task]
command = celery -A task worker
directory = /home/iems5722
user = root
autostart = true
autorestart = true
stdout_logfile = /home/iems5722/task.log
redirect_stderr = true